# launchd/com.youtube-channels.auto-post-comment.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.youtube-channels.auto-post-comment</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/mba/01-dev/youtube-channels/automation/auto_post_comment.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>StandardOutPath</key>
    <string>/Users/mba/01-dev/youtube-channels/logs/auto_post_comment_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mba/01-dev/youtube-channels/logs/auto_post_comment_stderr.log</string>
</dict>
</plist>
